# rootfs/home/_TEMP/.config/user/xmonad.target
[Unit]
Description=xmonad
After=graphical.target
Wants=graphical.target xorg.target
